# telemetry.service
[Service]
Type=exec
ExecStart=%h/dev/timeseries-lab/bin/telemetry
Restart=always

[Install]
WantedBy=default.target
